# zwave-js-ui/rootfs/etc/s6-overlay/s6-rc.d/init-zwave-js-ui/up
/etc/s6-overlay/s6-rc.d/init-zwave-js-ui/run
