# zwave-js-ui/rootfs/etc/services.d/zwave-js-ui/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Z-Wave JS UI
# Take down the S6 supervision tree when zwave-js-ui fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Z-Wave JS UI crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Z-Wave JS UI stoped, restarting..."
